# contrib/systemd/moxidle.service.in
[Install]
WantedBy=graphical-session.target

[Unit]
Description=moxidle idle manager
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStart=@bindir@/moxidle
Restart=always
RestartSec=10
